# deb/govee-gateway/DEBIAN/prerm
#!/bin/sh
# PRERM script for goveebttemplogger

echo "\033[36m Removing govee-gateway.\033[39m"
systemctl stop govee-gateway.service
systemctl disable govee-gateway.service
exit 0
